# target/linux/realtek/rtl8196e/base-files/etc/init.d/swdone
#!/bin/sh /etc/rc.common
START=99

start() {
	#reset bootting counter in hw info to indicate booting succesfully	
	echo 0 > /proc/flash/bootcnt
}


